# scheduled_task_persistence.yml
title: Suspicious Scheduled Task Persistence Executing Temp BAT
id: 0b88ac7d-3f5a-41ef-a3f0-schtasks-tempbat
status: experimental
description: Detects schtasks.exe creating scheduled tasks that execute .bat payloads from user Temp directory, commonly used for malware persistence.
author: k3rnelcallz
date: 2025/12/30
references:
  - https://attack.mitre.org/techniques/T1053/005/

logsource:
  product: windows
  category: process_creation

detection:
  selection_schtasks:
    Image|endswith: '\schtasks.exe'

  selection_create:
    CommandLine|contains:
      - '/create'
      - '/sc'
      - '/tn'
      - '/tr'

  selection_persistence_trigger:
    CommandLine|contains:
      - 'onlogon'
      - 'onstartup'
      - 'logon'
      - 'reboot'

  selection_temp_bat:
    CommandLine|contains:
      - '\AppData\Local\Temp'
      - '.bat'

  condition: selection_schtasks and selection_create and selection_temp_bat and selection_persistence_trigger

fields:
  - UtcTime
  - Image
  - CommandLine
  - ParentImage
  - User
  - IntegrityLevel
  - ParentCommandLine

falsepositives:
  - legitimate enterprise scheduled automation tasks
  - admin-created tasks for scripts

level: high

tags:
  - attack.persistence
  - attack.t1053.005
  - attack.execution